# root/etc/services.d/jackett/run
#!/command/execlineb -P

exec /command/s6-setuidgid jackett /opt/Jackett/jackett -d /config/jackett --NoUpdates